# data/mycdn-daemon.service
[Unit]
Description=MyCDN service
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/mycdn-daemon

[Install]
WantedBy=multi-user.target
